# com.jinwoo.anatomy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  매일 오전 8시에 send_daily.py 를 실행하는 macOS launchd 설정.
  설치:
    1) 아래 __PROJECT_PATH__ 두 곳을 실제 경로로 바꾸기 (예: /Users/jinwoo/anatomy-study)
    2) cp com.jinwoo.anatomy.plist ~/Library/LaunchAgents/
    3) launchctl load ~/Library/LaunchAgents/com.jinwoo.anatomy.plist
  해제: launchctl unload ~/Library/LaunchAgents/com.jinwoo.anatomy.plist
  시간 변경: 아래 Hour/Minute 값 수정 후 unload → load 다시.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.jinwoo.anatomy</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>__PROJECT_PATH__/send_daily.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__PROJECT_PATH__</string>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>8</integer>
    <key>Minute</key><integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>__PROJECT_PATH__/send.log</string>
  <key>StandardErrorPath</key>
  <string>__PROJECT_PATH__/send.err.log</string>
</dict>
</plist>
