# com.jinwoo.anatomy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  매일 오전 8시에 send_daily.py 를 실행하는 macOS launchd 설정.
  경로(/Users/jinwoo/anatomy-study)는 이미 채워져 있음.
  설치:
    1) cp com.jinwoo.anatomy.plist ~/Library/LaunchAgents/
    2) launchctl load ~/Library/LaunchAgents/com.jinwoo.anatomy.plist
  해제: launchctl unload ~/Library/LaunchAgents/com.jinwoo.anatomy.plist
  시간 변경: 아래 Hour/Minute 값 수정 후 unload → load 다시.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.jinwoo.anatomy</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>/Users/jinwoo/anatomy-study/send_daily.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/jinwoo/anatomy-study</string>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>8</integer>
    <key>Minute</key><integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/jinwoo/anatomy-study/send.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/jinwoo/anatomy-study/send.err.log</string>
</dict>
</plist>
